# example/cloudflare-helper.timer
[Unit]
Description=Cloudflare-Helper scheduled update

[Timer]
OnCalendar=Sat 04:00
Persistent=true

[Install]
WantedBy=timers.target